# supervisor/conf.d/ripd.conf
[program:ripd]
command=/usr/sbin/ripd -f /etc/quagga/ripd.conf
